# Repository/rc.local
#!/bin/sh
#/etc/scripts/interfaces/zombie0.sh &
#(sleep 20
/etc/scripts/mount/mount.sh) &
