# Repository/rc.local
#!/bin/sh
#/etc/scripts/interfaces/zombie0.sh &
