# Repository/rc.local
#!/bin/sh
#set -e
#/etc/scripts/interfaces/zombie0.sh &
#(sleep 10
#/etc/scripts/mount/mount.sh) &
#(sleep 30
#/etc/scripts/tunnels/strychnine.sh) &
